# fcopy.service
[Unit]
Description=FCopy daemon

[Service]
Type=simple
ExecStart=/usr/bin/hv_fcopy_daemon -n

[Install]
WantedBy=multi-user.target
