# vss.service
[Unit]
Description=VSS daemon

[Service]
Type=simple
ExecStart=/usr/bin/hv_vss_daemon -n

[Install]
WantedBy=multi-user.target
